# etc/udev/rules.d/99-sonos-bluetooth.rules
ACTION=="add|remove", KERNEL=="input[0-9]*", RUN+="/app/bin/sonos-bluetooth"
